# everule/cron/crontab
*/2      *       *       *       *       root *everulecron/download_upstream

